# csun_cat_and_sch.service
[Unit]
Description=CSUN Catalog and Schedules
After=network.target

[Service]
Type=forking
ExecStart=/bin/bash ~/CSUN-Catalog-And-Schedules/csun_cat_and_sch.sh start
ExecStop=/bin/bash ~/CSUN-Catalog-And-Schedules/csun_cat_and_sch.sh stop
ExecReload=/bin/bash ~/CSUN-Catalog-And-Schedules/csun_cat_and_sch.sh restart
User=kyeou

[Install]
WantedBy=multi-user.target